# scripts/vps-monitor.service
[Unit]
Description=VPS Monitor — System metrics HTTP server for Surge Panel
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/surge
Environment=VPS_MONITOR_TOKEN=changeme
ExecStart=/usr/bin/python3 /root/surge/scripts/vps_monitor.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
